# dev-env-examples/sshd-with-simple-http/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock
pidfile=/tmp/supervisord.pid
childlogdir=/tmp/supervisor
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[program:flask_app]
command=python3 app.py
directory=/app
autorestart=true
killasgroup=true
stopasgroup=true

; logs
stdout_logfile = /app/flask.log ; Where to write log messages
redirect_stderr=true
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=5

; restart
autostart=true
startretries=20
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
directory=/app
autorestart=true
killasgroup=true
stopasgroup=true

; logs
stdout_logfile = /app/sshd.log ; Where to write log messages
redirect_stderr=true
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=5

; restart
autostart=true
startretries=20
autorestart=true